# terraform/k8s-setup/templates/vault.service.tpl
[Unit]
Description = "Starts Vault"

[Service]
ExecStart = /usr/bin/vault agent -config=${vault_config}

[Install]
WantedBy = multi-user.target
